# etc/docker/aspnet/conf/etc/supervisord/apache2.conf
[program:apache2]
command=/opt/docker/bin/service/apache2.sh
numprocs=1
autostart=true
autorestart=true