# etc/docker/aspnet/conf/etc/supervisord/dotnet.conf
[program:dotnet]
command=/opt/docker/bin/service/dotnet.sh
numprocs=1
autostart=true
autorestart=true
